# system/snap-lxd-23991.mount
[Unit]
Description=Mount unit for lxd, revision 23991
Before=snapd.service
After=zfs-mount.service

[Mount]
What=/var/lib/snapd/snaps/lxd_23991.snap
Where=/snap/lxd/23991
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=default.target multi-user.target
